# deployment/systemd/klippyai-agent.service
[Unit]
Description=KlippyAI agent
After=network-online.target moonraker.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/KlippyAI
EnvironmentFile=-/etc/klippyai/klippyai.env
ExecStart=/home/pi/KlippyAI/.venv/bin/klippyai-agent
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
